# linux/handy-ptt.service
[Unit]
Description=Handy push-to-talk via foot pedal and/or F1 key

[Service]
ExecStart=/usr/bin/python3 %h/dotfiles/linux/handy-ptt.py
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
